# finance-dashboard/rootfs/etc/services.d/finance-dashboard/run
#!/usr/bin/with-contenv bashio

set -euo pipefail

# Supervisor's s6 service wrapper executes this script.
exec /usr/local/bin/run.sh
